# backend/deploy/gogomarket-backend.service
[Unit]
Description=GoGoMarket backend module (uploads + reviews)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/projects/video-marketplace-app
EnvironmentFile=/etc/gogomarket-backend.env
Environment=PYTHONUNBUFFERED=1

# Keep it simple: run on loopback, nginx proxies /api to 127.0.0.1:8000
ExecStart=/home/deploy/projects/video-marketplace-app/backend/.venv/bin/uvicorn backend.app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
